# supervisord_programs/rq.conf
[program:worker]
directory=/app
command=rq
    worker
    airscore-jobs
numprocs=%(ENV_RQ_WORKERS)s
stopsignal=TERM
autostart=true
autorestart=true

